# images.solidsense.io/images.solidsense.io.timer
[Unit]
Description=Update images.solidsense.io every 15m

[Timer]
OnCalendar=*-*-* *:*:15
Persistent=true

[Install]
WantedBy=timers.target
